# laravel-worker.conf
[program:mangala-ar]  
process_name=%(program_name)s_%(process_num)02d  
command=php /var/www/html/mangala-ar/artisan queue:listen --tries=3  
autostart=true  
autorestart=true  
stopasgroup=true  
killasgroup=true  
user=www-data
numprocs=1  
redirect_stderr=true  
stdout_logfile=/var/www/html/mangala-ar/storage/logs/worker.log  
stopwaitsecs=30
