# laravel-worker.conf
[program:m-ai]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/m-ai/artisan queue:listen --tries=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/html/m-ai/storage/logs/worker.log
stopwaitsecs=30
